# systemd/fiery-eyes-scan.timer
[Unit]
Description=Fiery Eyes — Hourly Scan Timer

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
