# systemd/fiery-eyes-report.service
[Unit]
Description=Fiery Eyes — Daily Nightly Strategist Report
After=network.target postgresql.service

[Service]
Type=oneshot
WorkingDirectory=/opt/fiery-eyes
ExecStart=/opt/fiery-eyes/venv/bin/python main.py report
EnvironmentFile=/opt/fiery-eyes/.env
User=fiery-eyes
StandardOutput=append:/opt/fiery-eyes/logs/report.log
StandardError=append:/opt/fiery-eyes/logs/report.log

[Install]
WantedBy=multi-user.target
